# deploy/habits-bot.service
# systemd-юнит для запуска Telegram-бота «Дневник привычек» на Linux-VPS 24/7.
#
# Установка (от root или через sudo):
#   1. Скопируй проект на сервер, например в /opt/habits-bot
#   2. Поправь User и пути ниже под свой сервер
#   3. sudo cp deploy/habits-bot.service /etc/systemd/system/
#   4. sudo systemctl daemon-reload
#   5. sudo systemctl enable --now habits-bot
#
# Логи смотреть:   journalctl -u habits-bot -f
# Перезапуск:      sudo systemctl restart habits-bot
# Остановка:       sudo systemctl stop habits-bot

[Unit]
Description=Telegram-бот «Дневник привычек»
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Пользователь, от которого работает бот (не root — создай отдельного, например: adduser --system habits)
User=habits
WorkingDirectory=/opt/habits-bot
# Запуск через python из виртуального окружения проекта
ExecStart=/opt/habits-bot/.venv/bin/python bot.py
# Перезапускать всегда при падении или штатном выходе
Restart=always
RestartSec=5
# Корректное логирование UTF-8 (на случай экзотичной локали сервера)
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONUTF8=1

[Install]
WantedBy=multi-user.target
